# examples/upstart/vnstat.conf
description "vnstat - network traffic monitor"

start on runlevel [2345]
stop on runlevel [016] or deconfiguring-networking or unmounting-filesystem

exec vnstatd -n
